# docker-pypi-mirror/files/etc/cron.d/bandersnatch
LC_ALL=en.US.utf8
0 0 * * * root bandersnatch mirror |& logger -t bandersnatch[mirror]
